# lesson6/spawn-cgi.socket
[Unit]
Description=Spawn FastCGI socket

[Socket]
ListenStream=/var/run/php-fcgi.sock

[Install]
WantedBy=sockets.target